# docker/supervisord/torspider.conf
[program:torspider]
command=python TorSpider.py